# include_resources/supervisord.conf

[program:tomcat]
command=/downloads/tomcat/bin/catalina.sh
autostart=true               
autorestart=true              
startsecs=10 
startretries=3 
